# infrastructure/scrimmage/crontab
00 00,12 * * * /bin/bash /app/matchmake.sh
